# swift/datanode/start_swift.conf
description "mount swift devices"
start on runlevel [234]
stop on runlevel [0156]
exec /usr/local/bin/mount_devices
